# totebot2.conf
# totebot2: The Quickening

start on runlevel [2345]
stop on runlevel [!2345]

respawn

source /etc/environment
script
    cd /home/ubuntu/totebot2/ && bin/hubot --adapter hipchat
end script