# totebot2.conf
# totebot2: The Quickening

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec source /etc/environment && cd /home/ubuntu/apps/totebot2/ && bin/hubot --adapter hipchat